# extra/gdm/.footprint
drwxr-xr-x	root/root	etc/
drwxr-xr-x	root/root	etc/X11/
drwxr-xr-x	root/root	etc/X11/dm/
drwxr-xr-x	root/root	etc/gdm/
drwxr-xr-x	root/root	etc/gdm/Init/
-rwxr-xr-x	root/root	etc/gdm/Init/Default
drwxr-xr-x	root/root	etc/gdm/PostLogin/
-rwxr-xr-x	root/root	etc/gdm/PostLogin/Default.sample
drwxr-xr-x	root/root	etc/gdm/PostSession/
-rwxr-xr-x	root/root	etc/gdm/PostSession/Default
drwxr-xr-x	root/root	etc/gdm/PreSession/
-rwxr-xr-x	root/root	etc/gdm/PreSession/Default
-rwxr-xr-x	root/root	etc/gdm/XKeepsCrashing
-rwxr-xr-x	root/root	etc/gdm/Xsession
-rw-r--r--	root/root	etc/gdm/custom.conf
-rw-r--r--	root/root	etc/gdm/custom.conf.orig
-rw-r--r--	root/root	etc/gdm/locale.alias
drwxr-xr-x	root/root	etc/gdm/modules/
-rw-r--r--	root/root	etc/gdm/modules/AccessDwellMouseEvents
-rw-r--r--	root/root	etc/gdm/modules/AccessKeyMouseEvents
-rw-r--r--	root/root	etc/gdm/modules/factory-AccessDwellMouseEvents
-rw-r--r--	root/root	etc/gdm/modules/factory-AccessKeyMouseEvents
drwxr-xr-x	root/root	etc/pam.d/
-rw-r--r--	root/root	etc/pam.d/gdm
-rw-r--r--	root/root	etc/pam.d/gdm-autologin
drwxr-xr-x	root/root	etc/rc.d/
drwxr-xr-x	root/root	etc/rc.d/init.d/
-rwxr-xr--	root/root	etc/rc.d/init.d/gdm
drwxr-xr-x	root/root	etc/rc.d/rc0.d/
lrwxrwxrwx	root/root	etc/rc.d/rc0.d/K05gdm -> ../init.d/gdm
drwxr-xr-x	root/root	etc/rc.d/rc1.d/
lrwxrwxrwx	root/root	etc/rc.d/rc1.d/K05gdm -> ../init.d/gdm
drwxr-xr-x	root/root	etc/rc.d/rc2.d/
lrwxrwxrwx	root/root	etc/rc.d/rc2.d/K05gdm -> ../init.d/gdm
drwxr-xr-x	root/root	etc/rc.d/rc3.d/
lrwxrwxrwx	root/root	etc/rc.d/rc3.d/K05gdm -> ../init.d/gdm
drwxr-xr-x	root/root	etc/rc.d/rc4.d/
lrwxrwxrwx	root/root	etc/rc.d/rc4.d/S95gdm -> ../init.d/gdm
drwxr-xr-x	root/root	etc/rc.d/rc5.d/
lrwxrwxrwx	root/root	etc/rc.d/rc5.d/K05gdm -> ../init.d/gdm
drwxr-xr-x	root/root	etc/rc.d/rc6.d/
lrwxrwxrwx	root/root	etc/rc.d/rc6.d/K05gdm -> ../init.d/gdm
drwxr-xr-x	root/root	etc/rc.d/rcsysinit.d/
drwxr-xr-x	root/root	etc/sysconfig/
drwxr-xr-x	root/root	usr/
drwxr-xr-x	root/root	usr/bin/
-rwxr-xr-x	root/root	usr/bin/gdm-dmx-reconnect-proxy
lrwxrwxrwx	root/root	usr/bin/gdmXnest -> gdmXnestchooser
-rwxr-xr-x	root/root	usr/bin/gdmXnestchooser
-rwxr-xr-x	root/root	usr/bin/gdmdynamic
-rwxr-xr-x	root/root	usr/bin/gdmflexiserver
-rwxr-xr-x	root/root	usr/bin/gdmphotosetup
-rwxr-xr-x	root/root	usr/bin/gdmthemetester
drwxr-xr-x	root/root	usr/lib/
-rwxr-xr-x	root/root	usr/lib/gdmaskpass
-rwxr-xr-x	root/root	usr/lib/gdmchooser
-rwxr-xr-x	root/root	usr/lib/gdmgreeter
-rwxr-xr-x	root/root	usr/lib/gdmlogin
-rwxr-xr-x	root/root	usr/lib/gdmopen
-rwxr-xr-x	root/root	usr/lib/gdmtranslate
drwxr-xr-x	root/root	usr/lib/gtk-2.0/
drwxr-xr-x	root/root	usr/lib/gtk-2.0/modules/
-rw-r--r--	root/root	usr/lib/gtk-2.0/modules/libdwellmouselistener.a
-rwxr-xr-x	root/root	usr/lib/gtk-2.0/modules/libdwellmouselistener.la
-rwxr-xr-x	root/root	usr/lib/gtk-2.0/modules/libdwellmouselistener.so
-rw-r--r--	root/root	usr/lib/gtk-2.0/modules/libkeymouselistener.a
-rwxr-xr-x	root/root	usr/lib/gtk-2.0/modules/libkeymouselistener.la
-rwxr-xr-x	root/root	usr/lib/gtk-2.0/modules/libkeymouselistener.so
drwxr-xr-x	root/root	usr/sbin/
-rwxr-xr-x	root/root	usr/sbin/gdm
-rwxr-xr-x	root/root	usr/sbin/gdm-binary
-rwxr-xr-x	root/root	usr/sbin/gdm-restart
-rwxr-xr-x	root/root	usr/sbin/gdm-safe-restart
-rwxr-xr-x	root/root	usr/sbin/gdm-stop
-rwxr-xr-x	root/root	usr/sbin/gdmsetup
drwxr-xr-x	root/root	usr/share/
drwxr-xr-x	root/root	usr/share/gdm/
drwxr-xr-x	root/root	usr/share/gdm/BuiltInSessions/
-rw-r--r--	root/root	usr/share/gdm/BuiltInSessions/default.desktop
drwxr-xr-x	root/root	usr/share/gdm/applications/
-rw-r--r--	root/root	usr/share/gdm/applications/gdmflexiserver-xnest.desktop
-rw-r--r--	root/root	usr/share/gdm/applications/gdmflexiserver.desktop
-rw-r--r--	root/root	usr/share/gdm/applications/gdmphotosetup.desktop
-rw-r--r--	root/root	usr/share/gdm/applications/gdmsetup.desktop
-r--r--r--	root/root	usr/share/gdm/defaults.conf
-rw-r--r--	root/root	usr/share/gdm/factory-defaults.conf
-rw-r--r--	root/root	usr/share/gdm/gdmchooser.glade
-rw-r--r--	root/root	usr/share/gdm/gdmphotosetup.glade
-rw-r--r--	root/root	usr/share/gdm/gdmsetup.glade
drwxr-xr-x	root/root	usr/share/gdm/themes/
drwxr-xr-x	root/root	usr/share/gdm/themes/circles/
-rw-r--r--	root/root	usr/share/gdm/themes/circles/GdmGreeterTheme.desktop
-rw-r--r--	root/root	usr/share/gdm/themes/circles/background.svg
-rw-r--r--	root/root	usr/share/gdm/themes/circles/circles.xml
-rw-r--r--	root/root	usr/share/gdm/themes/circles/flower.png
-rw-r--r--	root/root	usr/share/gdm/themes/circles/screenshot.png
drwxr-xr-x	root/root	usr/share/gdm/themes/happygnome-list/
-rw-r--r--	root/root	usr/share/gdm/themes/happygnome-list/GdmGreeterTheme.desktop
-rw-r--r--	root/root	usr/share/gdm/themes/happygnome-list/background.svg
-rw-r--r--	root/root	usr/share/gdm/themes/happygnome-list/gnome-logo.svg
-rw-r--r--	root/root	usr/share/gdm/themes/happygnome-list/happygnome-list.xml
-rw-r--r--	root/root	usr/share/gdm/themes/happygnome-list/screenshot.png
drwxr-xr-x	root/root	usr/share/gdm/themes/happygnome/
-rw-r--r--	root/root	usr/share/gdm/themes/happygnome/GdmGreeterTheme.desktop
-rw-r--r--	root/root	usr/share/gdm/themes/happygnome/background.svg
-rw-r--r--	root/root	usr/share/gdm/themes/happygnome/gnome-logo.svg
-rw-r--r--	root/root	usr/share/gdm/themes/happygnome/happygnome.xml
-rw-r--r--	root/root	usr/share/gdm/themes/happygnome/screenshot.png
drwxr-xr-x	root/root	usr/share/gdm/themes/nutyx09/
-rwxr-xr-x	install/install	usr/share/gdm/themes/nutyx09/Arc-GDM.xml
-rwxr-xr-x	install/install	usr/share/gdm/themes/nutyx09/GdmGreeterTheme.desktop
-rwxr-xr-x	install/install	usr/share/gdm/themes/nutyx09/bg.png
drwxr-xr-x	install/install	usr/share/gdm/themes/nutyx09/gtk-2.0/
-rwxr-xr-x	install/install	usr/share/gdm/themes/nutyx09/gtk-2.0/gtkrc
-rwxr-xr-x	install/install	usr/share/gdm/themes/nutyx09/halt.png
-rwxr-xr-x	install/install	usr/share/gdm/themes/nutyx09/main_box.png
-rwxr-xr-x	install/install	usr/share/gdm/themes/nutyx09/reboot.png
-rwxr-xr-x	install/install	usr/share/gdm/themes/nutyx09/screen.jpg
-rwxr-xr-x	install/install	usr/share/gdm/themes/nutyx09/suspend.png
-rwxr-xr-x	install/install	usr/share/gdm/themes/nutyx09/text.png
drwxr-xr-x	root/root	usr/share/gnome/
drwxr-xr-x	root/root	usr/share/gnome/help/
drwxr-xr-x	root/root	usr/share/gnome/help/gdm/
drwxr-xr-x	root/root	usr/share/gnome/help/gdm/C/
-rw-r--r--	root/root	usr/share/gnome/help/gdm/C/gdm.xml
-rw-r--r--	root/root	usr/share/gnome/help/gdm/C/legal.xml
drwxr-xr-x	root/root	usr/share/gnome/help/gdm/de/
-rw-r--r--	root/root	usr/share/gnome/help/gdm/de/gdm.xml
drwxr-xr-x	root/root	usr/share/gnome/help/gdm/en_GB/
-rw-r--r--	root/root	usr/share/gnome/help/gdm/en_GB/gdm.xml
drwxr-xr-x	root/root	usr/share/gnome/help/gdm/es/
-rw-r--r--	root/root	usr/share/gnome/help/gdm/es/gdm.xml
drwxr-xr-x	root/root	usr/share/gnome/help/gdm/fr/
-rw-r--r--	root/root	usr/share/gnome/help/gdm/fr/gdm.xml
drwxr-xr-x	root/root	usr/share/gnome/help/gdm/ko/
-rw-r--r--	root/root	usr/share/gnome/help/gdm/ko/gdm.xml
drwxr-xr-x	root/root	usr/share/gnome/help/gdm/oc/
-rw-r--r--	root/root	usr/share/gnome/help/gdm/oc/gdm.xml
drwxr-xr-x	root/root	usr/share/gnome/help/gdm/sv/
-rw-r--r--	root/root	usr/share/gnome/help/gdm/sv/gdm.xml
drwxr-xr-x	root/root	usr/share/gnome/help/gdm/uk/
-rw-r--r--	root/root	usr/share/gnome/help/gdm/uk/gdm.xml
drwxr-xr-x	root/root	usr/share/icons/
drwxr-xr-x	root/root	usr/share/icons/hicolor/
drwxr-xr-x	root/root	usr/share/icons/hicolor/16x16/
drwxr-xr-x	root/root	usr/share/icons/hicolor/16x16/apps/
-rw-r--r--	root/root	usr/share/icons/hicolor/16x16/apps/gdm-xnest.png
-rw-r--r--	root/root	usr/share/icons/hicolor/16x16/apps/gdmflexiserver.png
-rw-r--r--	root/root	usr/share/icons/hicolor/16x16/apps/gdmsetup.png
drwxr-xr-x	root/root	usr/share/icons/hicolor/22x22/
drwxr-xr-x	root/root	usr/share/icons/hicolor/22x22/apps/
-rw-r--r--	root/root	usr/share/icons/hicolor/22x22/apps/gdmflexiserver.png
-rw-r--r--	root/root	usr/share/icons/hicolor/22x22/apps/gdmsetup.png
drwxr-xr-x	root/root	usr/share/icons/hicolor/24x24/
drwxr-xr-x	root/root	usr/share/icons/hicolor/24x24/apps/
-rw-r--r--	root/root	usr/share/icons/hicolor/24x24/apps/gdmflexiserver.png
-rw-r--r--	root/root	usr/share/icons/hicolor/24x24/apps/gdmsetup.png
drwxr-xr-x	root/root	usr/share/icons/hicolor/32x32/
drwxr-xr-x	root/root	usr/share/icons/hicolor/32x32/apps/
-rw-r--r--	root/root	usr/share/icons/hicolor/32x32/apps/gdm-xnest.png
-rw-r--r--	root/root	usr/share/icons/hicolor/32x32/apps/gdmflexiserver.png
-rw-r--r--	root/root	usr/share/icons/hicolor/32x32/apps/gdmsetup.png
drwxr-xr-x	root/root	usr/share/icons/hicolor/48x48/
drwxr-xr-x	root/root	usr/share/icons/hicolor/48x48/apps/
-rw-r--r--	root/root	usr/share/icons/hicolor/48x48/apps/gdm-login-photo.png
-rw-r--r--	root/root	usr/share/icons/hicolor/48x48/apps/gdm-xnest.png
-rw-r--r--	root/root	usr/share/icons/hicolor/48x48/apps/gdm.png
-rw-r--r--	root/root	usr/share/icons/hicolor/48x48/apps/gdmflexiserver.png
-rw-r--r--	root/root	usr/share/icons/hicolor/48x48/apps/gdmsetup.png
drwxr-xr-x	root/root	usr/share/icons/hicolor/scalable/
drwxr-xr-x	root/root	usr/share/icons/hicolor/scalable/apps/
-rw-r--r--	root/root	usr/share/icons/hicolor/scalable/apps/gdmflexiserver.svg
-rw-r--r--	root/root	usr/share/icons/hicolor/scalable/apps/gdmsetup.svg
drwxr-xr-x	root/root	usr/share/man/
drwxr-xr-x	root/root	usr/share/man/man1/
-rw-r--r--	root/root	usr/share/man/man1/gdm.1.gz
drwxr-xr-x	root/root	usr/share/omf/
drwxr-xr-x	root/root	usr/share/omf/gdm/
-rw-r--r--	root/root	usr/share/omf/gdm/gdm-C.omf
-rw-r--r--	root/root	usr/share/omf/gdm/gdm-de.omf
-rw-r--r--	root/root	usr/share/omf/gdm/gdm-en_GB.omf
-rw-r--r--	root/root	usr/share/omf/gdm/gdm-es.omf
-rw-r--r--	root/root	usr/share/omf/gdm/gdm-fr.omf
-rw-r--r--	root/root	usr/share/omf/gdm/gdm-ko.omf
-rw-r--r--	root/root	usr/share/omf/gdm/gdm-oc.omf
-rw-r--r--	root/root	usr/share/omf/gdm/gdm-sv.omf
-rw-r--r--	root/root	usr/share/omf/gdm/gdm-uk.omf
drwxr-xr-x	root/root	usr/share/pixmaps/
drwxr-xr-x	root/root	usr/share/pixmaps/faces/
-rw-r--r--	root/root	usr/share/pixmaps/faces/astronaut.jpg
-rw-r--r--	root/root	usr/share/pixmaps/faces/baseball.png
-rw-r--r--	root/root	usr/share/pixmaps/faces/butterfly.png
-rw-r--r--	root/root	usr/share/pixmaps/faces/cat-eye.jpg
-rw-r--r--	root/root	usr/share/pixmaps/faces/chess.jpg
-rw-r--r--	root/root	usr/share/pixmaps/faces/coffee.jpg
-rw-r--r--	root/root	usr/share/pixmaps/faces/dice.jpg
-rw-r--r--	root/root	usr/share/pixmaps/faces/energy-arc.jpg
-rw-r--r--	root/root	usr/share/pixmaps/faces/fish.jpg
-rw-r--r--	root/root	usr/share/pixmaps/faces/flake.jpg
-rw-r--r--	root/root	usr/share/pixmaps/faces/flower.jpg
-rw-r--r--	root/root	usr/share/pixmaps/faces/grapes.jpg
-rw-r--r--	root/root	usr/share/pixmaps/faces/guitar.jpg
-rw-r--r--	root/root	usr/share/pixmaps/faces/launch.jpg
-rw-r--r--	root/root	usr/share/pixmaps/faces/leaf.jpg
-rw-r--r--	root/root	usr/share/pixmaps/faces/lightning.jpg
-rw-r--r--	root/root	usr/share/pixmaps/faces/penguin.jpg
-rw-r--r--	root/root	usr/share/pixmaps/faces/puppy.jpg
-rw-r--r--	root/root	usr/share/pixmaps/faces/sky.jpg
-rw-r--r--	root/root	usr/share/pixmaps/faces/soccerball.png
-rw-r--r--	root/root	usr/share/pixmaps/faces/sunflower.jpg
-rw-r--r--	root/root	usr/share/pixmaps/faces/sunset.jpg
-rw-r--r--	root/root	usr/share/pixmaps/faces/tennis-ball.png
-rw-r--r--	root/root	usr/share/pixmaps/faces/yellow-rose.jpg
-rw-r--r--	root/root	usr/share/pixmaps/gdm-foot-logo.png
-rw-r--r--	root/root	usr/share/pixmaps/nobody.png
-rw-r--r--	root/root	usr/share/pixmaps/nohost.png
drwxr-xr-x	root/root	usr/share/xsessions/
drwxr-xr-x	root/root	var/
drwxr-xr-x	root/root	var/lib/
drwxrwx--T	root/gdm	var/lib/gdm/
drwxr-xr-x	root/root	var/log/
drwxr-xr-x	root/root	var/log/gdm/
