# extra/networkmanager/.footprint
drwxr-xr-x	root/root	lib/
drwxr-xr-x	root/root	lib/udev/
drwxr-xr-x	root/root	lib/udev/rules.d/
-rw-r--r--	root/root	lib/udev/rules.d/77-nm-olpc-mesh.rules
drwxr-xr-x	root/root	usr/
drwxr-xr-x	root/root	usr/bin/
-rwxr-xr-x	root/root	usr/bin/nm-tool
drwxr-xr-x	root/root	usr/etc/
drwxr-xr-x	root/root	usr/etc/NetworkManager/
drwxr-xr-x	root/root	usr/etc/NetworkManager/dispatcher.d/
drwxr-xr-x	root/root	usr/etc/NetworkManager/system-connections/
drwxr-xr-x	root/root	usr/etc/dbus-1/
drwxr-xr-x	root/root	usr/etc/dbus-1/system.d/
-rw-r--r--	root/root	usr/etc/dbus-1/system.d/NetworkManager.conf
-rw-r--r--	root/root	usr/etc/dbus-1/system.d/nm-avahi-autoipd.conf
-rw-r--r--	root/root	usr/etc/dbus-1/system.d/nm-dhcp-client.conf
-rw-r--r--	root/root	usr/etc/dbus-1/system.d/nm-dispatcher.conf
-rw-r--r--	root/root	usr/etc/dbus-1/system.d/nm-ifcfg-rh.conf
drwxr-xr-x	root/root	usr/etc/rc.d/
drwxr-xr-x	root/root	usr/etc/rc.d/init.d/
-rwxr-xr-x	root/root	usr/etc/rc.d/init.d/NetworkManager
drwxr-xr-x	root/root	usr/include/
drwxr-xr-x	root/root	usr/include/NetworkManager/
-rw-r--r--	root/root	usr/include/NetworkManager/NetworkManager.h
-rw-r--r--	root/root	usr/include/NetworkManager/NetworkManagerVPN.h
-rw-r--r--	root/root	usr/include/NetworkManager/nm-connection.h
-rw-r--r--	root/root	usr/include/NetworkManager/nm-setting-8021x.h
-rw-r--r--	root/root	usr/include/NetworkManager/nm-setting-bluetooth.h
-rw-r--r--	root/root	usr/include/NetworkManager/nm-setting-cdma.h
-rw-r--r--	root/root	usr/include/NetworkManager/nm-setting-connection.h
-rw-r--r--	root/root	usr/include/NetworkManager/nm-setting-gsm.h
-rw-r--r--	root/root	usr/include/NetworkManager/nm-setting-ip4-config.h
-rw-r--r--	root/root	usr/include/NetworkManager/nm-setting-ip6-config.h
-rw-r--r--	root/root	usr/include/NetworkManager/nm-setting-olpc-mesh.h
-rw-r--r--	root/root	usr/include/NetworkManager/nm-setting-ppp.h
-rw-r--r--	root/root	usr/include/NetworkManager/nm-setting-pppoe.h
-rw-r--r--	root/root	usr/include/NetworkManager/nm-setting-serial.h
-rw-r--r--	root/root	usr/include/NetworkManager/nm-setting-vpn.h
-rw-r--r--	root/root	usr/include/NetworkManager/nm-setting-wired.h
-rw-r--r--	root/root	usr/include/NetworkManager/nm-setting-wireless-security.h
-rw-r--r--	root/root	usr/include/NetworkManager/nm-setting-wireless.h
-rw-r--r--	root/root	usr/include/NetworkManager/nm-setting.h
-rw-r--r--	root/root	usr/include/NetworkManager/nm-utils.h
drwxr-xr-x	root/root	usr/include/libnm-glib/
-rw-r--r--	root/root	usr/include/libnm-glib/libnm_glib.h
-rw-r--r--	root/root	usr/include/libnm-glib/nm-access-point.h
-rw-r--r--	root/root	usr/include/libnm-glib/nm-active-connection.h
-rw-r--r--	root/root	usr/include/libnm-glib/nm-cdma-device.h
-rw-r--r--	root/root	usr/include/libnm-glib/nm-client.h
-rw-r--r--	root/root	usr/include/libnm-glib/nm-device-bt.h
-rw-r--r--	root/root	usr/include/libnm-glib/nm-device-ethernet.h
-rw-r--r--	root/root	usr/include/libnm-glib/nm-device-wifi.h
-rw-r--r--	root/root	usr/include/libnm-glib/nm-device.h
-rw-r--r--	root/root	usr/include/libnm-glib/nm-dhcp4-config.h
-rw-r--r--	root/root	usr/include/libnm-glib/nm-exported-connection.h
-rw-r--r--	root/root	usr/include/libnm-glib/nm-gsm-device.h
-rw-r--r--	root/root	usr/include/libnm-glib/nm-ip4-config.h
-rw-r--r--	root/root	usr/include/libnm-glib/nm-ip6-config.h
-rw-r--r--	root/root	usr/include/libnm-glib/nm-object.h
-rw-r--r--	root/root	usr/include/libnm-glib/nm-remote-connection.h
-rw-r--r--	root/root	usr/include/libnm-glib/nm-remote-settings-system.h
-rw-r--r--	root/root	usr/include/libnm-glib/nm-remote-settings.h
-rw-r--r--	root/root	usr/include/libnm-glib/nm-serial-device.h
-rw-r--r--	root/root	usr/include/libnm-glib/nm-settings-connection-interface.h
-rw-r--r--	root/root	usr/include/libnm-glib/nm-settings-interface.h
-rw-r--r--	root/root	usr/include/libnm-glib/nm-settings-service.h
-rw-r--r--	root/root	usr/include/libnm-glib/nm-settings-system-interface.h
-rw-r--r--	root/root	usr/include/libnm-glib/nm-types.h
-rw-r--r--	root/root	usr/include/libnm-glib/nm-vpn-connection.h
-rw-r--r--	root/root	usr/include/libnm-glib/nm-vpn-plugin-ui-interface.h
-rw-r--r--	root/root	usr/include/libnm-glib/nm-vpn-plugin.h
drwxr-xr-x	root/root	usr/lib/
drwxr-xr-x	root/root	usr/lib/NetworkManager/
-rw-r--r--	root/root	usr/lib/NetworkManager/libnm-settings-plugin-ifcfg-rh.a
-rwxr-xr-x	root/root	usr/lib/NetworkManager/libnm-settings-plugin-ifcfg-rh.la
-rwxr-xr-x	root/root	usr/lib/NetworkManager/libnm-settings-plugin-ifcfg-rh.so
-rw-r--r--	root/root	usr/lib/NetworkManager/libnm-settings-plugin-keyfile.a
-rwxr-xr-x	root/root	usr/lib/NetworkManager/libnm-settings-plugin-keyfile.la
-rwxr-xr-x	root/root	usr/lib/NetworkManager/libnm-settings-plugin-keyfile.so
-rw-r--r--	root/root	usr/lib/libnm-glib-vpn.a
-rwxr-xr-x	root/root	usr/lib/libnm-glib-vpn.la
lrwxrwxrwx	root/root	usr/lib/libnm-glib-vpn.so -> libnm-glib-vpn.so.1.0.0
lrwxrwxrwx	root/root	usr/lib/libnm-glib-vpn.so.1 -> libnm-glib-vpn.so.1.0.0
-rwxr-xr-x	root/root	usr/lib/libnm-glib-vpn.so.1.0.0
-rw-r--r--	root/root	usr/lib/libnm-glib.a
-rwxr-xr-x	root/root	usr/lib/libnm-glib.la
lrwxrwxrwx	root/root	usr/lib/libnm-glib.so -> libnm-glib.so.2.2.0
lrwxrwxrwx	root/root	usr/lib/libnm-glib.so.2 -> libnm-glib.so.2.2.0
-rwxr-xr-x	root/root	usr/lib/libnm-glib.so.2.2.0
-rw-r--r--	root/root	usr/lib/libnm-util.a
-rwxr-xr-x	root/root	usr/lib/libnm-util.la
lrwxrwxrwx	root/root	usr/lib/libnm-util.so -> libnm-util.so.1.1.0
lrwxrwxrwx	root/root	usr/lib/libnm-util.so.1 -> libnm-util.so.1.1.0
-rwxr-xr-x	root/root	usr/lib/libnm-util.so.1.1.0
drwxr-xr-x	root/root	usr/lib/pkgconfig/
-rw-r--r--	root/root	usr/lib/pkgconfig/NetworkManager.pc
-rw-r--r--	root/root	usr/lib/pkgconfig/libnm-glib-vpn.pc
-rw-r--r--	root/root	usr/lib/pkgconfig/libnm-glib.pc
-rw-r--r--	root/root	usr/lib/pkgconfig/libnm-util.pc
drwxr-xr-x	root/root	usr/lib/pppd/
drwxr-xr-x	root/root	usr/lib/pppd/2.4.4/
-rw-r--r--	root/root	usr/lib/pppd/2.4.4/nm-pppd-plugin.a
-rwxr-xr-x	root/root	usr/lib/pppd/2.4.4/nm-pppd-plugin.la
-rwxr-xr-x	root/root	usr/lib/pppd/2.4.4/nm-pppd-plugin.so
drwxr-xr-x	root/root	usr/libexec/
-rwxr-xr-x	root/root	usr/libexec/nm-avahi-autoipd.action
-rwxr-xr-x	root/root	usr/libexec/nm-crash-logger
-rwxr-xr-x	root/root	usr/libexec/nm-dhcp-client.action
-rwxr-xr-x	root/root	usr/libexec/nm-dispatcher.action
drwxr-xr-x	root/root	usr/sbin/
-rwxr-xr-x	root/root	usr/sbin/NetworkManager
drwxr-xr-x	root/root	usr/share/
drwxr-xr-x	root/root	usr/share/NetworkManager/
-rw-r--r--	root/root	usr/share/NetworkManager/gdb-cmd
drwxr-xr-x	root/root	usr/share/dbus-1/
drwxr-xr-x	root/root	usr/share/dbus-1/system-services/
-rw-r--r--	root/root	usr/share/dbus-1/system-services/org.freedesktop.nm_dispatcher.service
drwxr-xr-x	root/root	usr/share/gtk-doc/
drwxr-xr-x	root/root	usr/share/gtk-doc/html/
drwxr-xr-x	root/root	usr/share/gtk-doc/html/libnm-glib/
-rw-r--r--	root/root	usr/share/gtk-doc/html/libnm-glib/NMAccessPoint.html
-rw-r--r--	root/root	usr/share/gtk-doc/html/libnm-glib/NMActiveConnection.html
-rw-r--r--	root/root	usr/share/gtk-doc/html/libnm-glib/NMCdmaDevice.html
-rw-r--r--	root/root	usr/share/gtk-doc/html/libnm-glib/NMClient.html
-rw-r--r--	root/root	usr/share/gtk-doc/html/libnm-glib/NMDHCP4Config.html
-rw-r--r--	root/root	usr/share/gtk-doc/html/libnm-glib/NMDevice.html
-rw-r--r--	root/root	usr/share/gtk-doc/html/libnm-glib/NMDeviceEthernet.html
-rw-r--r--	root/root	usr/share/gtk-doc/html/libnm-glib/NMDeviceWifi.html
-rw-r--r--	root/root	usr/share/gtk-doc/html/libnm-glib/NMGsmDevice.html
-rw-r--r--	root/root	usr/share/gtk-doc/html/libnm-glib/NMIP4Config.html
-rw-r--r--	root/root	usr/share/gtk-doc/html/libnm-glib/NMIP6Config.html
-rw-r--r--	root/root	usr/share/gtk-doc/html/libnm-glib/NMObject.html
-rw-r--r--	root/root	usr/share/gtk-doc/html/libnm-glib/NMSerialDevice.html
-rw-r--r--	root/root	usr/share/gtk-doc/html/libnm-glib/ch01.html
-rw-r--r--	root/root	usr/share/gtk-doc/html/libnm-glib/home.png
-rw-r--r--	root/root	usr/share/gtk-doc/html/libnm-glib/index.html
-rw-r--r--	root/root	usr/share/gtk-doc/html/libnm-glib/index.sgml
-rw-r--r--	root/root	usr/share/gtk-doc/html/libnm-glib/left.png
-rw-r--r--	root/root	usr/share/gtk-doc/html/libnm-glib/libnm-glib-NMDBusConnection.html
-rw-r--r--	root/root	usr/share/gtk-doc/html/libnm-glib/libnm-glib-NMDBusSettings.html
-rw-r--r--	root/root	usr/share/gtk-doc/html/libnm-glib/libnm-glib-NMDBusSettingsSystem.html
-rw-r--r--	root/root	usr/share/gtk-doc/html/libnm-glib/libnm-glib-NMSettings.html
-rw-r--r--	root/root	usr/share/gtk-doc/html/libnm-glib/libnm-glib-NMVPNConnection.html
-rw-r--r--	root/root	usr/share/gtk-doc/html/libnm-glib/libnm-glib-NMVPNPlugin.html
-rw-r--r--	root/root	usr/share/gtk-doc/html/libnm-glib/libnm-glib-libnm-glib.html
-rw-r--r--	root/root	usr/share/gtk-doc/html/libnm-glib/libnm-glib-nm-types.html
-rw-r--r--	root/root	usr/share/gtk-doc/html/libnm-glib/libnm-glib-nm-vpn-plugin-ui-interface.html
-rw-r--r--	root/root	usr/share/gtk-doc/html/libnm-glib/libnm-glib.devhelp
-rw-r--r--	root/root	usr/share/gtk-doc/html/libnm-glib/libnm-glib.devhelp2
-rw-r--r--	root/root	usr/share/gtk-doc/html/libnm-glib/right.png
-rw-r--r--	root/root	usr/share/gtk-doc/html/libnm-glib/style.css
-rw-r--r--	root/root	usr/share/gtk-doc/html/libnm-glib/up.png
drwxr-xr-x	root/root	usr/share/gtk-doc/html/libnm-util/
-rw-r--r--	root/root	usr/share/gtk-doc/html/libnm-util/NMConnection.html
-rw-r--r--	root/root	usr/share/gtk-doc/html/libnm-util/NMSetting.html
-rw-r--r--	root/root	usr/share/gtk-doc/html/libnm-util/NMSetting8021x.html
-rw-r--r--	root/root	usr/share/gtk-doc/html/libnm-util/NMSettingBluetooth.html
-rw-r--r--	root/root	usr/share/gtk-doc/html/libnm-util/NMSettingCdma.html
-rw-r--r--	root/root	usr/share/gtk-doc/html/libnm-util/NMSettingConnection.html
-rw-r--r--	root/root	usr/share/gtk-doc/html/libnm-util/NMSettingGsm.html
-rw-r--r--	root/root	usr/share/gtk-doc/html/libnm-util/NMSettingIP4Config.html
-rw-r--r--	root/root	usr/share/gtk-doc/html/libnm-util/NMSettingIP6Config.html
-rw-r--r--	root/root	usr/share/gtk-doc/html/libnm-util/NMSettingOlpcMesh.html
-rw-r--r--	root/root	usr/share/gtk-doc/html/libnm-util/NMSettingPPP.html
-rw-r--r--	root/root	usr/share/gtk-doc/html/libnm-util/NMSettingPPPOE.html
-rw-r--r--	root/root	usr/share/gtk-doc/html/libnm-util/NMSettingSerial.html
-rw-r--r--	root/root	usr/share/gtk-doc/html/libnm-util/NMSettingVPN.html
-rw-r--r--	root/root	usr/share/gtk-doc/html/libnm-util/NMSettingWired.html
-rw-r--r--	root/root	usr/share/gtk-doc/html/libnm-util/NMSettingWireless.html
-rw-r--r--	root/root	usr/share/gtk-doc/html/libnm-util/NMSettingWirelessSecurity.html
-rw-r--r--	root/root	usr/share/gtk-doc/html/libnm-util/api-index-full.html
-rw-r--r--	root/root	usr/share/gtk-doc/html/libnm-util/ch01.html
-rw-r--r--	root/root	usr/share/gtk-doc/html/libnm-util/home.png
-rw-r--r--	root/root	usr/share/gtk-doc/html/libnm-util/index.html
-rw-r--r--	root/root	usr/share/gtk-doc/html/libnm-util/index.sgml
-rw-r--r--	root/root	usr/share/gtk-doc/html/libnm-util/left.png
-rw-r--r--	root/root	usr/share/gtk-doc/html/libnm-util/libnm-util-nm-utils.html
-rw-r--r--	root/root	usr/share/gtk-doc/html/libnm-util/libnm-util.devhelp
-rw-r--r--	root/root	usr/share/gtk-doc/html/libnm-util/libnm-util.devhelp2
-rw-r--r--	root/root	usr/share/gtk-doc/html/libnm-util/object-tree.html
-rw-r--r--	root/root	usr/share/gtk-doc/html/libnm-util/right.png
-rw-r--r--	root/root	usr/share/gtk-doc/html/libnm-util/style.css
-rw-r--r--	root/root	usr/share/gtk-doc/html/libnm-util/up.png
drwxr-xr-x	root/root	usr/share/locale/
drwxr-xr-x	root/root	usr/share/locale/fr/
drwxr-xr-x	root/root	usr/share/locale/fr/LC_MESSAGES/
-rw-r--r--	root/root	usr/share/locale/fr/LC_MESSAGES/NetworkManager.mo
drwxr-xr-x	root/root	usr/share/man/
drwxr-xr-x	root/root	usr/share/man/man1/
-rw-r--r--	root/root	usr/share/man/man1/nm-tool.1.gz
drwxr-xr-x	root/root	usr/share/man/man8/
-rw-r--r--	root/root	usr/share/man/man8/NetworkManager.8.gz
drwxr-xr-x	root/root	usr/share/polkit-1/
drwxr-xr-x	root/root	usr/share/polkit-1/actions/
-rw-r--r--	root/root	usr/share/polkit-1/actions/org.freedesktop.network-manager-settings.system.policy
drwxr-xr-x	root/root	usr/var/
drwxr-xr-x	root/root	usr/var/run/
drwx------	root/root	usr/var/run/NetworkManager/
