# extra/ufw/.footprint
drwxr-xr-x	root/root	etc/
drwxr-xr-x	root/root	etc/default/
-rw-r--r--	root/root	etc/default/ufw
drwxr-xr-x	root/root	etc/rc.d/
drwxr-xr-x	root/root	etc/rc.d/init.d/
-rwxr-xr-x	root/root	etc/rc.d/init.d/ufwd
drwxr-xr-x	root/root	etc/rc.d/rc1.d/
drwxr-xr-x	root/root	etc/rc.d/rc2.d/
lrwxrwxrwx	root/root	etc/rc.d/rc2.d/S19ufwd -> ../init.d/ufwd
drwxr-xr-x	root/root	etc/rc.d/rc3.d/
lrwxrwxrwx	root/root	etc/rc.d/rc3.d/S19ufwd -> ../init.d/ufwd
drwxr-xr-x	root/root	etc/rc.d/rc4.d/
lrwxrwxrwx	root/root	etc/rc.d/rc4.d/S19ufwd -> ../init.d/ufwd
drwxr-xr-x	root/root	etc/rc.d/rc5.d/
lrwxrwxrwx	root/root	etc/rc.d/rc5.d/S19ufwd -> ../init.d/ufwd
drwxr-xr-x	root/root	etc/ufw/
-rw-r--r--	root/root	etc/ufw/after.rules
-rw-r--r--	root/root	etc/ufw/after6.rules
drwxr-xr-x	root/root	etc/ufw/applications.d/
-rw-r--r--	root/root	etc/ufw/before.rules
-rw-r--r--	root/root	etc/ufw/before6.rules
-rw-r--r--	root/root	etc/ufw/sysctl.conf
-rw-r--r--	root/root	etc/ufw/ufw.conf
drwxr-xr-x	root/root	lib/
drwxr-xr-x	root/root	lib/ufw/
-rwxr-xr-x	root/root	lib/ufw/ufw-init
-rwxr-xr-x	root/root	lib/ufw/ufw-init-functions
-rw-r--r--	root/root	lib/ufw/user.rules
-rw-r--r--	root/root	lib/ufw/user6.rules
drwxr-xr-x	root/root	usr/
drwxr-xr-x	root/root	usr/bin/
-rwxr-xr-x	root/root	usr/bin/ufw
drwxr-xr-x	root/root	usr/lib/
drwxr-xr-x	root/root	usr/lib/python2.6/
drwxr-xr-x	root/root	usr/lib/python2.6/site-packages/
-rw-r--r--	root/root	usr/lib/python2.6/site-packages/ufw-0.29.1-py2.6.egg-info
drwxr-xr-x	root/root	usr/lib/python2.6/site-packages/ufw/
-rw-r--r--	root/root	usr/lib/python2.6/site-packages/ufw/__init__.py (EMPTY)
-rw-r--r--	root/root	usr/lib/python2.6/site-packages/ufw/__init__.pyc
-rw-r--r--	root/root	usr/lib/python2.6/site-packages/ufw/applications.py
-rw-r--r--	root/root	usr/lib/python2.6/site-packages/ufw/applications.pyc
-rw-r--r--	root/root	usr/lib/python2.6/site-packages/ufw/backend.py
-rw-r--r--	root/root	usr/lib/python2.6/site-packages/ufw/backend.pyc
-rw-r--r--	root/root	usr/lib/python2.6/site-packages/ufw/backend_iptables.py
-rw-r--r--	root/root	usr/lib/python2.6/site-packages/ufw/backend_iptables.pyc
-rw-r--r--	root/root	usr/lib/python2.6/site-packages/ufw/common.py
-rw-r--r--	root/root	usr/lib/python2.6/site-packages/ufw/common.pyc
-rw-r--r--	root/root	usr/lib/python2.6/site-packages/ufw/frontend.py
-rw-r--r--	root/root	usr/lib/python2.6/site-packages/ufw/frontend.pyc
-rw-r--r--	root/root	usr/lib/python2.6/site-packages/ufw/util.py
-rw-r--r--	root/root	usr/lib/python2.6/site-packages/ufw/util.pyc
drwxr-xr-x	root/root	usr/share/
drwxr-xr-x	root/root	usr/share/man/
drwxr-xr-x	root/root	usr/share/man/man8/
-rw-r--r--	root/root	usr/share/man/man8/ufw-framework.8.gz
-rw-r--r--	root/root	usr/share/man/man8/ufw.8.gz
drwxr-xr-x	root/root	usr/share/ufw/
drwxr-xr-x	root/root	usr/share/ufw/messages/
-rw-r--r--	root/root	usr/share/ufw/messages/ar.mo
-rw-r--r--	root/root	usr/share/ufw/messages/bg.mo
-rw-r--r--	root/root	usr/share/ufw/messages/ca.mo
-rw-r--r--	root/root	usr/share/ufw/messages/cs.mo
-rw-r--r--	root/root	usr/share/ufw/messages/da.mo
-rw-r--r--	root/root	usr/share/ufw/messages/de.mo
-rw-r--r--	root/root	usr/share/ufw/messages/el.mo
-rw-r--r--	root/root	usr/share/ufw/messages/en_AU.mo
-rw-r--r--	root/root	usr/share/ufw/messages/en_GB.mo
-rw-r--r--	root/root	usr/share/ufw/messages/fi.mo
-rw-r--r--	root/root	usr/share/ufw/messages/fr.mo
-rw-r--r--	root/root	usr/share/ufw/messages/he.mo
-rw-r--r--	root/root	usr/share/ufw/messages/hu.mo
-rw-r--r--	root/root	usr/share/ufw/messages/id.mo
-rw-r--r--	root/root	usr/share/ufw/messages/it.mo
-rw-r--r--	root/root	usr/share/ufw/messages/nb.mo
-rw-r--r--	root/root	usr/share/ufw/messages/nl.mo
-rw-r--r--	root/root	usr/share/ufw/messages/pl.mo
-rw-r--r--	root/root	usr/share/ufw/messages/pt.mo
-rw-r--r--	root/root	usr/share/ufw/messages/pt_BR.mo
-rw-r--r--	root/root	usr/share/ufw/messages/ru.mo
-rw-r--r--	root/root	usr/share/ufw/messages/sk.mo
-rw-r--r--	root/root	usr/share/ufw/messages/sl.mo
-rw-r--r--	root/root	usr/share/ufw/messages/sr.mo
-rw-r--r--	root/root	usr/share/ufw/messages/sv.mo
-rw-r--r--	root/root	usr/share/ufw/messages/tl.mo
-rw-r--r--	root/root	usr/share/ufw/messages/zh_CN.mo
