# extra/slim/.footprint
drwxr-xr-x	root/root	etc/
drwxr-xr-x	root/root	etc/rc.d/
drwxr-xr-x	root/root	etc/rc.d/init.d/
-rwxr-xr-x	root/root	etc/rc.d/init.d/slim
drwxr-xr-x	root/root	etc/rc.d/rc0.d/
lrwxrwxrwx	root/root	etc/rc.d/rc0.d/K10slim -> ../init.d/slim
drwxr-xr-x	root/root	etc/rc.d/rc1.d/
lrwxrwxrwx	root/root	etc/rc.d/rc1.d/K10slim -> ../init.d/slim
drwxr-xr-x	root/root	etc/rc.d/rc2.d/
lrwxrwxrwx	root/root	etc/rc.d/rc2.d/K10slim -> ../init.d/slim
drwxr-xr-x	root/root	etc/rc.d/rc3.d/
lrwxrwxrwx	root/root	etc/rc.d/rc3.d/S40slim -> ../init.d/slim
drwxr-xr-x	root/root	etc/rc.d/rc4.d/
lrwxrwxrwx	root/root	etc/rc.d/rc4.d/K10slim -> ../init.d/slim
drwxr-xr-x	root/root	etc/rc.d/rc5.d/
lrwxrwxrwx	root/root	etc/rc.d/rc5.d/K10slim -> ../init.d/slim
drwxr-xr-x	root/root	etc/rc.d/rc6.d/
lrwxrwxrwx	root/root	etc/rc.d/rc6.d/K10slim -> ../init.d/slim
drwxr-xr-x	root/root	etc/skel/
-rw-r--r--	root/root	etc/slim.conf
-rw-r--r--	root/root	etc/slim.xinitrc
drwxr-xr-x	root/root	usr/
drwxr-xr-x	root/root	usr/bin/
-rwxr-xr-x	root/root	usr/bin/slim
drwxr-xr-x	root/root	usr/man/
drwxr-xr-x	root/root	usr/man/man1/
-rw-r--r--	root/root	usr/man/man1/slim.1.gz
drwxr-xr-x	root/root	usr/share/
drwxr-xr-x	root/root	usr/share/slim/
drwxr-xr-x	root/root	usr/share/slim/themes/
drwxr-xr-x	root/root	usr/share/slim/themes/default/
-rw-r--r--	root/root	usr/share/slim/themes/default/background.jpg
-rw-r--r--	root/root	usr/share/slim/themes/default/panel.png
-rw-r--r--	root/root	usr/share/slim/themes/default/slim.theme
drwxr-xr-x	root/root	usr/share/slim/themes/nutyx/
-rw-r--r--	root/root	usr/share/slim/themes/nutyx/background.jpg
-rw-r--r--	root/root	usr/share/slim/themes/nutyx/panel.png
-rw-r--r--	root/root	usr/share/slim/themes/nutyx/slim.theme
